# Detections/lolbins.yml
title: LOLRansom Rundll32 Suspicious Execution
id: win_lolransom_rundll32_suspicious
status: experimental
description: Detects abuse of rundll32.exe for LOLRansom-style fileless execution via mshtml/RunHTMLApplication and non-standard DLL entrypoints.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_base:
    Image|endswith: '\rundll32.exe'
  sel_cmd:
    CommandLine|contains:
      - 'javascript:'
      - 'RunHTMLApplication'
      - '.dll,'
      - 'mshtml'
  sel_parent:
    ParentImage|endswith:
      - '\explorer.exe'
      - '\winword.exe'
      - '\excel.exe'
      - '\outlook.exe'
  condition: sel_base and sel_cmd and sel_parent
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - ParentCommandLine
  - User
  - Hashes
falsepositives:
  - Rare legitimate rundll32 calls with custom exports (should be allowlisted per-app).
level: high
tags:
  - attack.t1218.011
  - attack.defense_evasion
  - attack.execution
---
title: LOLRansom Regsvr32 Remote Scriptlet
id: win_lolransom_regsvr32_remote_sct
status: experimental
description: Detects regsvr32.exe loading remote .sct scriptlets or scrobj.dll with URL-based /i parameter.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_img:
    Image|endswith: '\regsvr32.exe'
  sel_cmd:
    CommandLine|contains:
      - '.sct'
      - 'scrobj.dll'
      - '/i:http'
      - '/i:https'
  condition: sel_img and sel_cmd
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - User
falsepositives:
  - Very uncommon; investigate all hits.
level: critical
tags:
  - attack.t1218.010
  - attack.defense_evasion
  - attack.execution
---
title: LOLRansom MSHTA Remote or LOLBin-Orchestrated Execution
id: win_lolransom_mshta_remote
status: experimental
description: Detects mshta.exe executing remote HTA content or spawning secondary interpreters such as powershell/cmd.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_img:
    Image|endswith: '\mshta.exe'
  sel_cmd:
    CommandLine|contains:
      - 'http://'
      - 'https://'
      - '.hta'
      - 'powershell'
      - 'cmd.exe'
  condition: sel_img and sel_cmd
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - User
falsepositives:
  - Legacy admin tools or installers using HTA (should be allowlisted).
level: high
tags:
  - attack.t1218.005
  - attack.defense_evasion
  - attack.execution
---
title: PowerShell Obfuscated Download-Execute Patterns
id: win_lolransom_powershell_obfuscated_download
status: experimental
description: Flags common loader primitives used in fileless stages (DownloadData, IEX, Base64 payloads).
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: powershell
detection:
  sel_evt:
    EventID: 4104
  sel_str:
    ScriptBlockText|contains:
      - 'Invoke-Expression'
      - 'DownloadData'
      - 'FromBase64String'
  sel_b64:
    ScriptBlockText|re: '(?i)[A-Za-z0-9+/]{200,}={0,2}'
  condition: sel_evt and (sel_str or sel_b64)
fields:
  - UtcTime
  - ScriptBlockText
  - User
  - HostApplication
falsepositives:
  - Admin scripts that legitimately decode data; validate by source and parent process.
level: high
tags:
  - attack.t1059.001
  - attack.execution
  - attack.defense_evasion
---
title: Registry Run/RunOnce Persistence via LOLBins
id: win_lolransom_registry_persistence
status: stable
description: Detects persistence creation using Run/RunOnce keys to launch mshta/regsvr32/powershell on logon.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_evt:
    EventID: 13
  sel_key:
    TargetObject|contains:
      - '\Software\Microsoft\Windows\CurrentVersion\Run'
      - '\Software\Microsoft\Windows\CurrentVersion\RunOnce'
  sel_val:
    Details|contains:
      - 'mshta.exe'
      - 'regsvr32.exe'
      - 'powershell.exe'
  condition: sel_evt and sel_key and sel_val
fields:
  - UtcTime
  - TargetObject
  - Details
  - Image
  - User
falsepositives:
  - Legitimate app updaters using these keys with signed binaries.
level: high
tags:
  - attack.t1547.001
  - attack.persistence
---
title: LOLRansom Reflective Injection Indicators via Process Access
id: win_lolransom_injection_access
status: experimental
description: Identifies suspicious process memory access that may precede reflective DLL injection or in-memory execution.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_evt:
    EventID: 10
  sel_src:
    SourceImage|endswith:
      - '\rundll32.exe'
      - '\regsvr32.exe'
      - '\mshta.exe'
      - '\powershell.exe'
  sel_mask:
    GrantedAccess|contains:
      - '0x1F0FFF'   # All access
      - '0x0010'     # VM_OPERATION
      - '0x0020'     # VM_READ
      - '0x0008'     # VM_WRITE
  condition: sel_evt and sel_src and sel_mask
fields:
  - UtcTime
  - SourceImage
  - TargetImage
  - GrantedAccess
  - CallTrace
  - User
falsepositives:
  - Security tools, debuggers. Build allowlist for known-good publishers.
level: high
tags:
  - attack.t1055
  - attack.defense_evasion
  - attack.privilege_escalation
---
title: LOLRansom Remote Thread Creation from LOLBins
id: win_lolransom_createremotethread
status: experimental
description: Detects CreateRemoteThread patterns often used post-write for code execution in another process.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_evt:
    EventID: 8
  sel_src:
    SourceImage|endswith:
      - '\rundll32.exe'
      - '\regsvr32.exe'
      - '\mshta.exe'
      - '\powershell.exe'
  condition: sel_evt and sel_src
fields:
  - UtcTime
  - SourceImage
  - TargetImage
  - StartAddress
  - StartModule
  - User
falsepositives:
  - Legit automation tools occasionally inject for control; verify signatures and publishers.
level: high
tags:
  - attack.t1055.001
  - attack.defense_evasion
---
title: Shadow Copy Deletion for Impact
id: win_lolransom_inhibit_recovery
status: stable
description: Detects vssadmin/wmic usage to delete shadow copies, a common ransomware preparatory action.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_img:
    Image|endswith:
      - '\vssadmin.exe'
      - '\wmic.exe'
  sel_cmd_vss:
    CommandLine|contains:
      - 'delete shadows'
      - 'resize shadowstorage'
  sel_cmd_wmic:
    CommandLine|contains:
      - 'shadowcopy delete'
  condition: sel_img and (sel_cmd_vss or sel_cmd_wmic)
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - User
falsepositives:
  - Backup/IT maintenance tasks (schedule known-good).
level: high
tags:
  - attack.t1490
  - attack.impact
---
title: Mass Extension Changes Suggestive of Encryption
id: win_lolransom_mass_rename
status: experimental
description: Heuristic to detect rapid creation/rename of files with uncommon extensions by suspicious processes.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_evt:
    EventID: 11
  sel_proc:
    Image|endswith:
      - '\rundll32.exe'
      - '\mshta.exe'
      - '\regsvr32.exe'
      - '\powershell.exe'
  sel_ext:
    TargetFilename|re: '\.(lock|crypt|enc|encrypted|.*\d{4})$'
  condition: sel_evt and sel_proc and sel_ext
fields:
  - UtcTime
  - Image
  - TargetFilename
  - ProcessGuid
  - ParentImage
  - User
falsepositives:
  - Bulk archival or compression jobs; correlate with user-initiated tasks.
level: high
tags:
  - attack.t1486
  - attack.impact
---
title: Zeek HTTP Suspicious Artifact Retrieval
id: net_lolransom_zeek_http_artifacts
status: experimental
description: Detects retrieval of common LOLRansom delivery artifacts (.hta/.sct) over HTTP(S) in Zeek logs.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: zeek
  service: http
detection:
  sel_uri:
    uri|endswith:
      - '.hta'
      - '.sct'
  sel_host:
    host|re: '^(?!.*(corp|intra|microsoft|windowsupdate|officecdn)\.).*$'
  condition: sel_uri and sel_host
fields:
  - ts
  - id.orig_h
  - id.resp_h
  - host
  - uri
  - user_agent
falsepositives:
  - Rare; developer/testing environments hosting HTA/SCT may trigger.
level: high
tags:
  - attack.t1071.001
  - attack.command_and_control
  - attack.resource_development
---
title: Suricata EVE Suspicious HTA/SCT Signature
id: net_lolransom_suricata_hta_sct
status: experimental
description: Matches Suricata alerts or HTTP logs indicating HTA/SCT downloads typical of LOLRansom staging.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: suricata
  service: eve
detection:
  sel_alert:
    event_type: alert
    alert.signature|contains:
      - '.hta'
      - '.sct'
  sel_http:
    event_type: http
    http.uri|endswith:
      - '.hta'
      - '.sct'
  condition: sel_alert or sel_http
fields:
  - timestamp
  - src_ip
  - src_port
  - dest_ip
  - dest_port
  - alert.signature
  - http.hostname
  - http.uri
falsepositives:
  - Low; validate hosting origin and business justification.
level: high
tags:
  - attack.t1071.001
  - attack.command_and_control
---
title: Parent-Child LOLBin Anomaly
id: win_lolransom_parent_child_anomaly
status: experimental
description: Detects unusual parent-child relationships where user-facing apps spawn LOLBins used by LOLRansom.
author: Copilot
date: 2025/08/28
references:
  - https://github.com/shenniyuki-coder/LOLRansom
logsource:
  product: windows
  service: sysmon
detection:
  sel_child:
    Image|endswith:
      - '\rundll32.exe'
      - '\regsvr32.exe'
      - '\mshta.exe'
      - '\powershell.exe'
  sel_parent:
    ParentImage|endswith:
      - '\winword.exe'
      - '\excel.exe'
      - '\powerpnt.exe'
      - '\outlook.exe'
      - '\acrord32.exe'
      - '\chrome.exe'
      - '\msedge.exe'
  condition: sel_child and sel_parent
fields:
  - UtcTime
  - Image
  - CommandLine
  - ParentImage
  - ParentCommandLine
  - User
falsepositives:
  - Add approved automation/macros to allowlist.
level: high
tags:
  - attack.t1204
  - attack.execution
  - attack.defense_evasion
